# deploy/systemd/sloptimize-api.service
[Unit]
Description=Sloptimize API Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/admin/sloptimize
Environment=PATH=/home/admin/sloptimize/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/home/admin/sloptimize/.env
ExecStart=/home/admin/sloptimize/.venv/bin/uvicorn sloptimize.api.main:app --host 0.0.0.0 --port 9090
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target